# scripts/launchd/com.example.meta-agent.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>de.example.meta-agent.api</string>
  <key>ProgramArguments</key><array><string>__META_AGENT_REPO__/scripts/supervisor/api.sh</string></array>
  <key>WorkingDirectory</key><string>__META_AGENT_REPO__</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>__META_AGENT_REPO__/logs/supervisor/api.out.log</string>
  <key>StandardErrorPath</key><string>__META_AGENT_REPO__/logs/supervisor/api.err.log</string>
</dict></plist>
